# etc/buildbot/slaves/openbsd-5.6-x86_64/crontab
MASTER_HOST=192.168.1.100
MASTER_PORT=9989
SLAVE_NAME=openbsd-5.6-x86_64
SLAVE_PASSWORD=pass
@reboot /usr/local/bin/buildslave restart /home/buildbot/libc11
