# etc/buildbot/slaves/netbsd-6.1.5-x86_64/crontab
MASTER_HOST=192.168.1.100
MASTER_PORT=9989
SLAVE_NAME=netbsd-6.1.5-x86_64
SLAVE_PASSWORD=pass
@reboot /usr/pkg/bin/buildslave restart /home/buildbot/libc11
