# etc/buildbot/master/crontab
SLAVE_PASSWORD=pass
@reboot /usr/local/bin/buildbot restart /home/buildmaster/libc11
