# etc/buildbot/slaves/ubuntu-14.04-armhf/crontab
MASTER_HOST=build.libc11.org
MASTER_PORT=9989
SLAVE_NAME=ubuntu-14.04-armhf
SLAVE_PASSWORD=pass
@reboot /usr/local/bin/buildslave restart /home/buildbot/libc11
